# deploy/factory-api.service
# systemd unit for the FastAPI backend (analytics API).
# Install: sudo cp deploy/factory-api.service /etc/systemd/system/
#          sudo systemctl daemon-reload && sudo systemctl enable --now factory-api
# Hetzner VPS layout: /home/scott/manufacturing-intelligence-platform.
# Ports 8000/8001/8787 are already used by other portfolio apps; API uses 8002.

[Unit]
Description=Manufacturing Intelligence API (FastAPI)
After=network.target postgresql.service

[Service]
Type=simple
User=scott
Group=scott
WorkingDirectory=/home/scott/manufacturing-intelligence-platform/backend
EnvironmentFile=/home/scott/manufacturing-intelligence-platform/backend/.env
Environment="PYTHONUNBUFFERED=1"
# Use a venv created with: python3 -m venv .venv && .venv/bin/pip install -r requirements.txt
ExecStart=/home/scott/manufacturing-intelligence-platform/backend/.venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8002 --workers 2
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
